# wework_chan.ini
[program:wework_chan]
directory = /root
command=/root/wework_chan -corp <替换为你的企业ID> -token <替换为你的token> -key <替换为你的encodingAESKey>
autostart = true
numprocs=1
user=root
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisord/wework_chan.log
loglevel=info